# musicadet-hud.service
[Unit]
Description=MusicaDet HUD - web dashboard
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/musicadet
ExecStart=/opt/musicadet/run-hud.sh
Restart=on-failure
RestartSec=5
User=root
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
