# musicadet.service
[Unit]
Description=MusicaDet - Spotify playlists and artist discographies via spotDL
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
WorkingDirectory=/opt/musicadet
ExecStart=/usr/bin/python3 /opt/musicadet/music_sync.py
User=root
Nice=10
IOSchedulingClass=best-effort
IOSchedulingPriority=7
StandardOutput=journal
StandardError=journal
TimeoutStartSec=21600
