# musicadet-hud.service
[Unit]
Description=MusicaDet HUD - web dashboard
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/musicadet
ExecStart=/usr/bin/python3 -m uvicorn hud:app --host 0.0.0.0 --port 8800
Restart=on-failure
RestartSec=5
User=root
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
